# roles/kubernetes_backup/templates/backup.timer
[Unit]
Description=Kubernetes Backup Timer

[Timer]
OnCalendar={{ kubernetes_backup_schedule }}

[Install]
WantedBy=timers.target
